# docs/deploy/djehuti-api.service
[Unit]
Description=Djehuti API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/djehuti-api
ExecStart=/opt/djehuti-api/Djehuti.Api
Restart=on-failure
RestartSec=5
Environment=ASPNETCORE_URLS=http://localhost:5008
Environment=ASPNETCORE_ENVIRONMENT=Production
# Add your Anthropic API key here:
# Environment=ANTHROPIC_API_KEY=sk-ant-...

[Install]
WantedBy=multi-user.target
